# crontab
0 14 * * * ubuntu /bin/bash /home/ubuntu/apps/bernard/repository/run_on_server.sh fab $DEPLOYMENT_TARGET cron_jobs.post_message >> /var/log/bernard/crontab.log 2>&1
